# s3fs.conf
[program:s3fs]
priority = 1
command = bash -c '/mount-s3.sh'
redirect_stderr = true
stdout_logfile = /var/log/secure
stdout_events_enabled = true
